# apps/raspberry-pi/config/pommai.service
[Unit]
Description=Pommai Smart Toy Client
After=network-online.target sound.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=10
User=pommai
Group=pommai
WorkingDirectory=/home/pommai/app
Environment="PATH=/home/pommai/app/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/home/pommai/app"
ExecStart=/home/pommai/app/venv/bin/python /home/pommai/app/pommai_client_fastrtc.py

# Resource limits for Pi Zero 2W
MemoryMax=200M
CPUQuota=60%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pommai

# Security hardening (relaxed for GPIO/I2C/SPI access)
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/home/pommai /tmp /var/log/pommai

[Install]
WantedBy=multi-user.target
